# puppet/modules/wiki/templates/irc-rc-relay.service.erb
# This file is managed by puppet
[Unit]
Description=IRC relay for MediaWiki recent changes
After=network.target

[Service]
User=<%= @maintenance_user %>
ExecStart=/usr/bin/perl <%= @config_dir %>/relays/irc/relay.pl
Restart=always
RestartSec=60
PrivateDevices=true
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=full
EnvironmentFile=/etc/irc-relay/irc-rc-relay.env

[Install]
WantedBy=multi-user.target
