# puppet/modules/wiki/templates/irc-log-relay.service.erb
# This file is managed by puppet
[Unit]
Description=IRC relay for PHP error log
After=network.target

[Service]
User=<%= @maintenance_user %>
ExecStart=/usr/bin/perl <%= @config_dir %>/irc-relay/bot/relay.pl
Restart=always
RestartSec=60
PrivateDevices=true
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=full

[Install]
WantedBy=multi-user.target
