# puppet/modules/wiki/templates/mw-jobrunner.service.erb
# This file is managed by puppet
[Unit]
Description=Process MediaWiki jobs
After=mysql.service network.target

[Service]
User=www-data
Group=www-data
ExecStart=/usr/bin/php <%= @production_dir %>/maintenance/run.php runJobs --wait --maxjobs=1000 --memory-limit=1G --procs=4
Restart=always
RestartSec=1
ProtectHome=read-only
ProtectSystem=full
Nice=5

[Install]
WantedBy=multi-user.target
