# puppet/modules/wiki/templates/mw-supportedlanguages.service.erb
# This file is managed by puppet
[Unit]
Description=Updated Special:SupportedLanguages caches
Requisite=mysql.service
After=mysql.service

[Service]
User=<%= @maintenance_user %>
ExecStart=/usr/bin/php <%= @production_dir %>/extensions/Translate/scripts/updateTranslatorActivity.php
ProtectHome=read-only
ProtectSystem=full
