# puppet/modules/wiki/templates/phplog2slack.service.erb
# This file is managed by puppet
[Unit]
Description=Process and relay error messages to slack-relay
After=network.target slack-relay.service
Requires=slack-relay.service

[Service]
User=<%= @maintenance_user %>
ExecStart=/bin/sh -c 'php <%= @config_dir %>/relays/logfilter.php > /run/slack-relay/pipe'
Restart=always
RestartSec=60
PrivateDevices=true
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=full

[Install]
WantedBy=multi-user.target
