# puppet/modules/awstats/templates/awstats.service.erb
# This file is managed by puppet
[Unit]
Description=Process logs for awstats

[Service]
ExecStart=/usr/lib/cgi-bin/awstats.pl -update -config=translatewiki.net
ProtectHome=read-only
ProtectSystem=full
