# puppet/modules/awstats/templates/mlebstats.service.erb
# This file is managed by puppet
[Unit]
Description=Emails MLEB download statistics

[Service]
ExecStart=/usr/bin/php /home/betawiki/config/bin/scripts/mlebstats.php /var/lib/awstats nlaxstrom@wikimedia.org,federico@translatewiki.net
ProtectHome=read-only
ProtectSystem=full
Nice=10
