# puppet/modules/awstats/files/awstats.cron
# file managed by puppet
0,10,20,30,40,50 * * * * root /usr/lib/cgi-bin/awstats.pl -update -config=translatewiki.net >/dev/null
0 6 1 * * root /usr/bin/php /home/betawiki/config/bin/scripts/mlebstats.php /var/lib/awstats nlaxstrom@wikimedia.org,federico@translatewiki.net
