# puppet/modules/repong/templates/autoexport.service.erb
# This file is managed by puppet
[Unit]
Description=Automatically exports non-MediaWiki translations
Requisite=mysql.service
After=mysql.service

[Service]
User=<%= @l10nbot_user %>
ExecStart=<%= @bin_dir %>/execute-and-notify <%= @bin_dir %>/<%= @command %>
ProtectHome=read-only
ProtectSystem=full
