# puppet/modules/repong/templates/autoimport.timer.erb
# This file is managed by puppet
[Unit]
Description=Run autoimport on a defined schedule

[Timer]
OnCalendar=<%= @when %>

[Install]
WantedBy=timers.target
