# puppet/modules/repong/templates/autoexport.timer.erb
# This file is managed by puppet
[Unit]
Description=Run autoexport on a defined schedule

[Timer]
# Our cron jobs for autoimport are running half past even hours (CET/CEST)
OnCalendar=<%= @when %>

[Install]
WantedBy=timers.target
