# puppet/modules/repong/templates/repong-cleanups.service.erb
# This file is managed by puppet
[Unit]
Description=Cleanups stale files related to repository handling

[Service]
ExecStart=<%= @bin_dir %>/repong-cleanups
ProtectHome=read-only
ProtectSystem=full
Nice=10
